# home/.config/autostart/chromium-browser.desktop
[Desktop Entry]
Exec=chromium-browser
Name=chromium-browser
Type=Application
