# home/.config/autostart/cirosantilli-autostart.desktop
[Desktop Entry]
Exec=cirosantilli-autostart
Name=cirosantilli-autostart
Type=Application
X-GNOME-Autostart-enabled=true
