# home/.config/autostart/firefox.desktop
[Desktop Entry]
Exec=firefox
Name=firefox
Type=Application
X-GNOME-Autostart-enabled=false
