# home/.config/autostart/firefox.desktop
[Desktop Entry]
Exec=myautostart
Name=myautostart
Type=Application
