# home/.config/autostart/gnome-terminal.desktop
[Desktop Entry]
Exec=gnome-terminal --window --full-screen -e tmux
Name=gnome-terminal
Type=Application
X-GNOME-Autostart-enabled=true
