# home/.config/autostart/guake.desktop
[Desktop Entry]
Exec=guake
Name=guake
Type=Application
X-GNOME-Autostart-enabled=false
