# home/.config/autostart/gvim.desktop
[Desktop Entry]
Exec=gvim
Name=gvim
Type=Application
