# home/.config/autostart/krusader.desktop
[Desktop Entry]
Exec=krusader
Name=krusader
Type=Application
X-GNOME-Autostart-enabled=false
