# home/.config/autostart/my-autostart.desktop
[Desktop Entry]
Exec=my-autostart
Name=my-autostart
Type=Application
X-GNOME-Autostart-enabled=true
