# home/.config/autostart/pidgin.desktop
[Desktop Entry]
Exec=pidgin
Name=pidgin
Type=Application
