# home/.config/autostart/skype.desktop
[Desktop Entry]
Exec=skype
Name=skype
Type=Application
