# home/.config/autostart/terminator.desktop
[Desktop Entry]
Exec=terminator -f -e tmux
Name=terminator
Type=Application
X-GNOME-Autostart-enabled=false
