# home/.config/autostart/terminator.desktop
[Desktop Entry]
Exec=terminator
Name=terminator
Type=Application
X-GNOME-Autostart-enabled=true
